# ops/beemuu-api.service
[Unit]
Description=BeeEmUu VPS backend API
After=network.target
Documentation=https://github.com/ohgeeceee/beemuu

[Service]
Type=simple
# Run as a module so relative imports in backend/ work. We set PYTHONPATH
# explicitly because systemd can strip cwd from sys.path under some conditions,
# and we hit a real-world `ImportError: attempted relative import with no known
# parent package` crashloop on the VPS before this fix.
ExecStart=/usr/bin/python3 -m backend.app --host 127.0.0.1 --port 8765
WorkingDirectory=/root/beemuu
Environment=PYTHONPATH=/root/beemuu
EnvironmentFile=/etc/beemuu/beemuu.env
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beemuu-api
User=root
Group=root

[Install]
WantedBy=multi-user.target
